# systemd/marker-sync.timer
[Unit]
Description=Run marker-sync every 2 minutes

[Timer]
OnBootSec=30s
OnUnitActiveSec=2min
RandomizedDelaySec=10s
Persistent=true
Unit=marker-sync.service

[Install]
WantedBy=timers.target